# ops/systemd/pg-backup.service
[Unit]
Description=IzdajIznajmi PostgreSQL backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
EnvironmentFile=/etc/izdaji/pg-backup.env
WorkingDirectory=/var/www/izdaji
ExecStart=/var/www/izdaji/ops/backup_pg.sh
