# ops/systemd/pg-backup.timer
[Unit]
Description=Run IzdajIznajmi PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
